# templates/default/systemd/ignite.socket.erb
#
# This file is generated by the Chef Cookbook: 
# Name:   Ignite
# Source: https://github.com/BegleyBrothers/cookbook-ignite
#
# Any changes will be overwritten.
#
[Unit]
Description=Ignite Socket for the daemon
PartOf=<%= @ignite_name %>.service

[Socket]
ListenStream=<%= @ignite_socket %>
SocketMode=0660
SocketUser=root
SocketGroup=<%= @config.group %>

[Install]
WantedBy=sockets.target
